# make-openwrt/openwrt-files/platform-files/amlogic/rootfs/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
ttyAML0::askfirst:/usr/libexec/login.sh
tty0::askfirst:/usr/libexec/login.sh
